# debian/DEBIAN/prerm
#!/bin/bash
systemctl stop amd-config-manager.service

systemctl disable amd-config-manager.service
